# etc/profile.d/activemq.sh
export JAVA_HOME=/usr/lib/jvm/jre
export PATH=$PATH:/opt/activemq/bin
export ACTIVEMQ_SSL_OPTS="-Djavax.net.ssl.keyStore=/opt/activemq/conf/broker.ks -Djavax.net.ssl.keyStorePassword=password"
